# conf/command.conf
[program:command]
command=bash /var/www/techstop/env.sh
